# server/murica-generator.service
[Unit]
Description=Run murica-generator.py

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/murica-generator.py
Environment="OPENAI_API_KEY=your_openai_api_key"